# nebulastreams.service
[Unit]
Description=NebulaStreams
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/sohil/NebulaStreams
Environment=HOME=/home/sohil
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PUBLIC_BASE_URL=https://nebulastreams.tail3ba349.ts.net
Environment=STREMIO_RESULT_MEMORY_CACHE_MAX_ENTRIES=200
Environment=STREMIO_RESULT_MEMORY_CACHE_MAX_MB=4
Environment=PROVIDER_RESULT_MEMORY_CACHE_MAX_ENTRIES=200
Environment=PROVIDER_RESULT_MEMORY_CACHE_MAX_MB=2
Environment=HUBCLOUD_MEMORY_CACHE_MAX_ENTRIES=100
Environment=HUBCLOUD_MEMORY_CACHE_MAX_MB=1
Environment=STREMIO_MAX_INFLIGHT_SEARCHES=35
Environment=MAX_ACTIVE_STREAMS=6
Environment=PROVIDER_GLOBAL_MAX_INFLIGHT=6
ExecStart=/home/sohil/NebulaStreams/start.sh
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
